# fan-control.timer
[Unit]
Description=Timer for fan-control

[Timer]
OnActiveSec=10
OnUnitActiveSec=10
Unit=fan-control.service

[Install]
WantedBy=timers.target
