# fan-control.service
[Unit]
Description=Run fan-control script

[Service]
Type=oneshot
ExecStart=/bin/sh /mnt/local/fancontrol.sh
